# com.betbot.mlb-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.betbot.mlb-update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/josephbarkie/Coding/betbot/machine_learning/scripts/schedule_updates.sh</string>
    </array>

    <!-- Run daily at 6:00 AM. launchd catches up missed runs on next wake,
         unlike cron which silently skips jobs fired while the machine sleeps. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/josephbarkie/Coding/betbot/logs/mlb_data_update.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/josephbarkie/Coding/betbot/logs/mlb_data_update.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
